# etc/cron.daily/meet-keys
#! /bin/bash

(apt-key list | grep "release@heroku.com") || apt-key add /etc/MEET/heroku.key
