# etc/cron.daily/meet-config-update
#!/bin/bash

#
# Update configuration files.
# [maxg]
#

( cd / && git pull --quiet origin master )

glib-compile-schemas /usr/share/glib-2.0/schemas
