# etc/cron.daily/meet-epoptes
#!/bin/bash

#
# Download Epoptes certificate.
# [maxg]
#

[[ -s /etc/epoptes/server.crt ]] || epoptes-client -c server
