# etc/network/if-up.d/epoptes-client
#!/bin/sh

# [MEET]
# Do not start as root
